# .config/autostart/RescueTime.desktop
[Desktop Entry]
Name=RescueTime
Comment=
Exec=rescuetime
Terminal=false
Type=Application