# .config/autostart/Set icons to papirus.desktop
[Desktop Entry]
Name=Set icons to papirus
Comment=Fix try icon issues
Exec=/home/adam/dev/scripts/iconset.sh
Terminal=false
Type=Application