# .config/autostart/Set icons to correct themes.desktop
[Desktop Entry]
Name=Set icons to correct themes
Comment=
Exec=iconset.sh
Terminal=false
Type=Application