# .config/autostart/Set up wtf config.desktop
[Desktop Entry]
Name=Set up wtf config
Comment=
Exec=/home/adam/dev/scripts/setup-wtf.sh
Terminal=false
Type=Application