# .config/autostart/Start NordVPN.desktop
[Desktop Entry]
Name=Start NordVPN
Comment=Run Start NordVPN
Exec=nordvpn connect US
Terminal=false
Type=Application